# shell/systemd/maginkcal.timer
[Unit]
Description=Run Magic Ink Calendar every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target