# shell/systemd/maginkcal.timer
[Unit]
Description=Run Magic Ink Calendar every hour

[Timer]
OnCalendar=*-*-* 06:00:00
OnCalendar=*-*-* 17:00:00

[Install]
WantedBy=timers.target